# vitodeploy/rootfs/etc/services.d/vito/run
#!/usr/bin/with-contenv bashio
set -euo pipefail

bashio::log.info "Starting Apache with Vito..."

# Apply config
export APP_KEY=$(bashio::config 'app_key')
export NAME=$(bashio::config 'name')
export EMAIL=$(bashio::config 'email')
export PASSWORD=$(bashio::config 'password')
export APP_URL=$(bashio::config 'app_url')

# Write environment variables for PHP
cat <<EOF >/var/www/html/.env
APP_KEY=${APP_KEY}
NAME=${NAME}
EMAIL=${EMAIL}
PASSWORD=${PASSWORD}
APP_URL=${APP_URL}
EOF

# Start Apache in foreground (pid 1 handled by HA's s6)
exec apache2ctl -D FOREGROUND
