# vitodeploy/rootfs/etc/s6-overlay/s6-rc.d/vito/finish
#!/usr/bin/with-contenv bashio

bashio::log.info "Vito service stopped"
