# vitodeploy/rootfs/etc/s6-overlay/s6-rc.d/vito/run
#!/usr/bin/with-contenv bashio

# Generate .env from add-on options
cat <<EOF >/var/www/html/.env
APP_KEY=$(bashio::config 'app_key')
NAME=$(bashio::config 'name')
EMAIL=$(bashio::config 'email')
PASSWORD=$(bashio::config 'password')
APP_URL=$(bashio::config 'app_url')
EOF

# Configure Apache for PHP
cat <<EOF >/etc/apache2/conf.d/vito.conf
<Directory "/var/www/html">
    Options Indexes FollowSymLinks
    AllowOverride All
    Require all granted
</Directory>

# PHP configuration
LoadModule php_module /usr/lib/apache2/mod_php83.so
AddType application/x-httpd-php .php
DirectoryIndex index.php index.html
EOF

bashio::log.info "Starting Vito web server..."

# Run Apache in foreground
exec httpd -D FOREGROUND -f /etc/apache2/httpd.conf
